# sv/wmname/up
#!/bin/execlineb -P
/usr/bin/s6-envdir /etc/s6/sv/wmname/env
wmname compiz
