# sv/slstatus/run
#!/bin/execlineb -P
s6-envdir env
slstatus
