# sv/xmodmap/up
#!/bin/execlineb -P
xmodmap -display :0 /etc/X11/Xmodmap
